# deploy/systemd/pmketoan-web.service
[Unit]
Description=Visota ERP Uvicorn Web Server
After=network.target mariadb.service

[Service]
Type=notify
User=visota
Group=www-data
WorkingDirectory=/home/visota/app
Environment="PATH=/home/visota/app/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings.prod"
EnvironmentFile=/home/visota/app/.env
ExecStart=/home/visota/app/.venv/bin/uvicorn config.asgi:application --host 127.0.0.1 --port 8000 --workers 4 --log-level info --access-log
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/visota /home/visota/app/staticfiles /home/visota/app/media
RestrictSUIDSGID=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target
