# deploy/systemd/visota-qcluster.service
[Unit]
Description=Visota ERP django-q2 Cluster (background workers)
After=network.target mariadb.service visota-web.service

[Service]
Type=simple
User=visota
Group=www-data
WorkingDirectory=/home/visota/app
Environment="PATH=/home/visota/app/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings.prod"
EnvironmentFile=/home/visota/app/.env
ExecStart=/home/visota/app/.venv/bin/python manage.py qcluster
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/visota /home/visota/app/media

[Install]
WantedBy=multi-user.target
