# deploy/systemd/visota-web.service
[Unit]
Description=Visota ERP Gunicorn Web Server
After=network.target mariadb.service

[Service]
Type=notify
User=visota
Group=www-data
WorkingDirectory=/home/visota/app
Environment="PATH=/home/visota/app/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings.prod"
EnvironmentFile=/home/visota/app/.env
ExecStart=/home/visota/app/.venv/bin/gunicorn --workers 4 --threads 2 --bind 127.0.0.1:8000 --access-logfile /var/log/visota/gunicorn-access.log --error-logfile /var/log/visota/gunicorn-error.log --timeout 120 --graceful-timeout 30 --max-requests 1000 --max-requests-jitter 100 config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/visota /home/visota/app/staticfiles /home/visota/app/media
RestrictSUIDSGID=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target
